# supervisor/etc/supervisor/conf.d/watchdog_graftnoded.conf
[program:graftnoded]
directory=/home/graft-sn/supernode/graft-sn-watchdog/
command=bash -c 'sleep 100 && ./gn.sh --log-level 1 > watchdog_gn.log 2>&1'
priority=20
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_backups=4
stdout_logfile_maxbytes=10MB
stderr_logfile=NONE
autostart=true
autorestart=true
startretries=0
numprocs=1